# packages/avrcp-volume/avrcp_volume-0.7.5-py3-none-any.whl/avrcp_volume/resources/avrcp-volume.service
[Unit]
Description=Avrcp volume controller
Requires=dbus.service bluetooth.target

[Service]
Type=simple
ExecStart=avrcp-volume

[Install]
WantedBy=user-graphical.target
