# packages/mqtt-to-influx/mqtt-to-influx-0.3.5.tar.gz/mqtt-to-influx-0.3.5/systemd/mqtt-to-influx.service
[Unit]
Description=MQTT To Influx
Requires=
After=network.target

[Service]
Type=simple
;EnvironmentFile=%h/configs/environment
ExecStart=mqtt-to-influx
;--config=%h/configs/%i.json --debug
Restart=on-failure
RestartSec=1m
User=mqttinflux

[Install]
WantedBy=default.target
