# packages/mpris-fakeplayer/mpris-fakeplayer-0.1.3.tar.gz/mpris-fakeplayer-0.1.3/mpris_fakeplayer/resources/mpris-fakeplayer.service
[Unit]
Description=Fake mpris player for activating bluez avrcp volume control
Requires=mpris-proxy.service
After=mpris-proxy.service

[Service]
Type=simple
ExecStart=mpris-fakeplayer

[Install]
WantedBy=default.target
